# scripts/com.efficiency-tracker.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.efficiency-tracker.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON_PATH__</string>
        <string>__PROJECT_PATH__/main.py</string>
        <string>--period</string>
        <string>day</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_PATH__</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>21</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__PROJECT_PATH__/logs/daily.log</string>

    <key>StandardErrorPath</key>
    <string>__PROJECT_PATH__/logs/daily.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
